# python/download sharepoint/queue2.conf
[program:queue-job2]

process_name=%(program_name)s_%(process_num)02d
command=sudo php /home/ilensys/compliance-qa/artisan queue:work --queue=file_download database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=4
redirect_stderr=true
stdout_logfile=/etc/supervisor/worker.log
stopwaitsecs=3600